# profile/airootfs/etc/systemd/system/sunglasses.service
[Unit]
Description=SunGlasses Session

[Service]
Type=simple
ExecStart=/home/sunny/config/apply.sh

[Install]
WantedBy=multi-user.target
